# systemd/vehicleagent.service
[Unit]
Description=RoadRunner vehicle telemetry service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/vehicleagent --simulator --mqtt-host homeassistant.local --base-topic car/pathfinder
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
